# mcp-server.service
[Unit]
Description=MCP Server - Custom tools for Claude Code
After=[TARGET PATH]
Wants=[TARGET PATH]

[Service]
Type=simple
User=[USER]
WorkingDirectory=/YOUR/PATH/TO/mcp-server
EnvironmentFile=/YOUR/PATH/TO/mcp-server/.env
ExecStart=/usr/bin/npx tsx src/index.ts
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/YOUR/PATH/
PrivateTmp=true

[Install]
WantedBy=multi-user.target
